# 00-k480.rules
KERNEL=="hidraw2", SUBSYSTEM=="hidraw", ACTION=="add", RUN+="/opt/bin/k480_conf -d /dev/hidraw2 -f on"
