# 00-k480.rules
SUBSYSTEM=="hidraw", ACTION=="add", RUN+="/opt/bin/k480.sh"
